# conf/prd/cityhallmonitor.cron
PATH=/home/apps/env/cityhallmonitor/bin:$PATH
MAILTO=knightlab@northwestern.edu
DJANGO_SETTINGS_MODULE=core.settings.prd

05 0,8,16 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data Person
10 0,8,16 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data BodyType
15 0,8,16 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data Body
20 0,8,16 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data MatterType
25 0,8,16 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data MatterStatus
30 0,8,16 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data Matter
30 1,9,17 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_sponsors
30 2,10,18 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_attachments
30 3,11,19 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_pdfs
30 4,12,20 * * * apps /home/apps/sites/cityhallmonitor/manage.py get_descriptions

