# conf/prd/cityhallmonitor.cron
PATH=/home/apps/env/cityhallmonitor/bin:$PATH
MAILTO=knightlab@northwestern.edu
DJANGO_SETTINGS_MODULE=core.settings.prd

05 * * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data Person >> /home/apps/log/cityhallmonitor/person.log
10 * * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data BodyType >> /home/apps/log/cityhallmonitor/bodytype.log
15 * * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data Body >> /home/apps/log/cityhallmonitor/body.log
20 * * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data MatterType >> /home/apps/log/cityhallmonitor/mattertype.log
25 * * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data MatterStatus >> /home/apps/log/cityhallmonitor/matterstatus.log
30 0,6,12,18 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data Matter >> /home/apps/log/cityhallmonitor/matter.log
30 1,7,13,19 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_attachments >> /home/apps/log/cityhallmonitor/matter_attachments.log
30 3,9,15,21 * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_pdfs >> /home/apps/log/cityhallmonitor/pdfs.log
30 5,11,17,23 * * * apps /home/apps/sites/cityhallmonitor/manage.py get_descriptions >> /home/apps/log/cityhallmonitor/descriptions.log
