# conf/prd/cityhallmonitor.cron
PATH=/home/apps/env/cityhallmonitor/bin:$PATH
MAILTO=knightlab@northwestern.edu
DJANGO_SETTINGS_MODULE=core.settings.prd

01 * * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_data Matter
10 * * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_attachments
25 * * * * apps /home/apps/sites/cityhallmonitor/manage.py pull_pdfs
50 * * * * apps /home/apps/sites/cityhallmonitor/manage.py get_descriptions
